# supervisord/conf.d/mysql.conf
[program:myslq]
command=service mysql start
